# tools/telegram-bot/agent-telegram-bot.service
[Unit]
Description=Agent Telegram Bot Service
Documentation=https://github.com/agent-setup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/etc/hive
ExecStart=/etc/hive/venv/bin/python3 /etc/hive/agent_telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/etc/hive /tmp
PrivateTmp=false

[Install]
WantedBy=multi-user.target
